# zzaqi.service
[Unit]
Description=zz_aqi

[Service]
Type=simple
ExecStart=/usr/local/bin/zz_aqi
