# zzaqi.timer
[Unit]
Description=Runs zz_aqi every hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h

[Install]
WantedBy=multi-user.target
